# mohomi/deploy/service/mihomo/run
#!/bin/sh
exec /usr/bin/mihomo -d /etc/mihomo